# GUI/remote.service
[Unit]
Description=Touchscreen Universal Remote
After=display-manager.service
Requires=display-manager.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 -u /opt/universal-remote/remote.py
WorkingDirectory=/opt/universal-remote
User=USERNAME

# Set env for X11; adjust DISPLAY if your session uses something else
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/USERNAME/.Xauthority
Environment=PYTHONUNBUFFERED=1
Restart=on-failure
RestartSec=5

# Optional: also write to files you can tail
StandardOutput=append:/opt/universal-remote/remote.out
StandardError=append:/opt/universal-remote/remote.err

[Install]
WantedBy=graphical.target
